# .config/systemd/user/borgbackup.service
[Service]
Type=oneshot
Environment=TARGET=mbx:/volatile/backup/borg
ExecStart=%h/scripts/borgbackup.sh
