# .config/systemd/user/borgbackup.service
[Service]
Type=oneshot
Environment=TARGET=mbx:%h/mnt/five/borg
ExecStart=%h/scripts/borgbackup.sh
