# .config/systemd/user/emacs.socket
[Socket]
ListenStream=/tmp/emacs%U/server

[Install]
WantedBy=sockets.target
