# .config/systemd/user/emacs-apps.socket
[Socket]
ListenStream=/run/user/%U/emacs/apps

[Install]
WantedBy=sockets.target
