# .config/systemd/user/kwin-restart.service
[Unit]
Description=restart kwin

[Service]
Type=oneshot
ExecStart=/bin/systemctl --user restart kwin
