# .config/systemd/user/kwin-restart.timer
[Unit]
Description=Rest kwin daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
