# .config/systemd/user/mediabox-backup.timer
[Unit]
Description=Backup mediabox

[Timer]
OnBootSec=15min
OnUnitActiveSec=1D

[Install]
WantedBy=timers.target
