# .config/systemd/user/sys-update.service
[Unit]
Description=Sys update

[Service]
Type=oneshot
ExecStart=%h/.config/systemd/scripts/ubu-update

[Install]
WantedBy=user.target
