# .config/systemd/user/trakt-scrobbler.service
[Unit]
Description=Trakt Scrobbler Service

[Service]
ExecStart="/home/fra/.local/bin/trakts" run

[Install]
WantedBy=default.target