# .config/etc/udev/70-snap.anbox.rules
# This file is automatically generated.
# udev tagging/device cgroups disabled with non-strict mode snaps
## opengl
#KERNEL=="nvhost-*", TAG+="snap_anbox_anbox"
## opengl
#KERNEL=="nvmap", TAG+="snap_anbox_anbox"
## opengl
#KERNEL=="renderD[0-9]*", TAG+="snap_anbox_anbox"
## opengl
#KERNEL=="tegra_dc_[0-9]*", TAG+="snap_anbox_anbox"
## opengl
#KERNEL=="tegra_dc_ctrl", TAG+="snap_anbox_anbox"
## opengl
#KERNEL=="vchiq", TAG+="snap_anbox_anbox"
## opengl
#SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_anbox_anbox"
#TAG=="snap_anbox_anbox", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_anbox_anbox $devpath $major:$minor"
## opengl
#KERNEL=="nvhost-*", TAG+="snap_anbox_android-settings"
## opengl
#KERNEL=="nvmap", TAG+="snap_anbox_android-settings"
## opengl
#KERNEL=="renderD[0-9]*", TAG+="snap_anbox_android-settings"
## opengl
#KERNEL=="tegra_dc_[0-9]*", TAG+="snap_anbox_android-settings"
## opengl
#KERNEL=="tegra_dc_ctrl", TAG+="snap_anbox_android-settings"
## opengl
#KERNEL=="vchiq", TAG+="snap_anbox_android-settings"
## opengl
#SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_anbox_android-settings"
#TAG=="snap_anbox_android-settings", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_anbox_android-settings $devpath $major:$minor"
## opengl
#KERNEL=="nvhost-*", TAG+="snap_anbox_appmgr"
## opengl
#KERNEL=="nvmap", TAG+="snap_anbox_appmgr"
## opengl
#KERNEL=="renderD[0-9]*", TAG+="snap_anbox_appmgr"
## opengl
#KERNEL=="tegra_dc_[0-9]*", TAG+="snap_anbox_appmgr"
## opengl
#KERNEL=="tegra_dc_ctrl", TAG+="snap_anbox_appmgr"
## opengl
#KERNEL=="vchiq", TAG+="snap_anbox_appmgr"
## opengl
#SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_anbox_appmgr"
#TAG=="snap_anbox_appmgr", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_anbox_appmgr $devpath $major:$minor"
