# .config/systemd/sys-update.service
[Unit]
Description=Sys update

[Service]
Type=oneshot
ExecStart=%h/.config/systemd/scripts/ubu-update
