# .config/systemd/user/ssh-agent.service.d/override.conf
[Service]
ExecStartPost=%h/.config/systemd/scripts/ssh/set-agent-pid.sh
